# systemd/nsfw-bot.service
[Unit]
Description=Telegram NSFW Detection Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nsfw-bot
Environment="PATH=/opt/nsfw-bot/.venv/bin"
ExecStart=/opt/nsfw-bot/.venv/bin/python bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nsfw-bot

[Install]
WantedBy=multi-user.target
