# launchd/com.telegram.nsfwbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.telegram.nsfwbot</string>
  <key>WorkingDirectory</key>
  <string>/Users/nishkarshkr/Desktop/bot-app</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nishkarshkr/Desktop/bot-app/.venv/bin/python</string>
    <string>/Users/nishkarshkr/Desktop/bot-app/bot.py</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/nishkarshkr/Library/Logs/telegram-nsfwbot.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nishkarshkr/Library/Logs/telegram-nsfwbot.err.log</string>
</dict>
</plist>